# deploy/netwatchdog.service
[Unit]
Description=netwatchdog — network port change monitor
Documentation=https://github.com/krudolf55/netwatchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=netwatchdog
Group=netwatchdog

ExecStartPre=/usr/local/bin/netwatchdog --config /etc/netwatchdog/netwatchdog.yaml init-db
ExecStart=/usr/local/bin/netwatchdog --config /etc/netwatchdog/netwatchdog.yaml start

Restart=on-failure
RestartSec=10s

# Directories created automatically by systemd (owner: netwatchdog)
StateDirectory=netwatchdog
LogsDirectory=netwatchdog
ConfigurationDirectory=netwatchdog

# Grant only the capabilities needed for raw packet scanning.
# Remove CAP_NET_ADMIN if you don't use masscan.
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN

# Harden the service process
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/netwatchdog /var/log/netwatchdog

# Allow nmap/masscan to open raw sockets
LockPersonality=yes
RestrictNamespaces=yes
RestrictRealtime=yes
SystemCallFilter=@system-service @network-io
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
